# etc/autostart/dummy.desktop
[Desktop Entry]
Name=Dummy
Icon=preferences-system-power-manager
Exec=dummy
Terminal=false
Type=Application
StartupNotify=true
