# etc/autostart/mountd.desktop
[Desktop Entry]
Name=Mount Daemon
Icon=preferences-system-power-manager
Exec=mountctl -p
Terminal=false
Type=Application
StartupNotify=true
