# etc/autostart/powerd.desktop
[Desktop Entry]
Name=Power Daemon
Icon=preferences-system-power-manager
Exec=powerctl -p
Terminal=false
Type=Application
StartupNotify=true
