# etc/autostart/powertray.desktop
[Desktop Entry]
Name=Power controller
Icon=preferences-system-power-manager
Exec=powertray
Terminal=false
Type=Application
StartupNotify=true
