# src/spmd/spmd.init
#!/sbin/runscript

depend() {
    use dns net
}

start() {
    ebegin "Starting Source Package Manager daemon"
    start-stop-daemon --start --quiet \
        --exec /usr/sbin/spmd
    eend $?
}

stop() {
    ebegin "Stopping Source Package Manager daemon"
    start-stop-daemon --stop --quiet \
        --exec /usr/sbin/spmd
    eend $?
}
